# pkg/deploy/templates/systemd.service.tmpl
[Unit]
Description={{APP_NAME}}
After=network.target

[Service]
WorkingDirectory=/srv/{{APP_NAME}}/current
EnvironmentFile=-/srv/{{APP_NAME}}/shared/env/.env
Environment=PORT={{PORT}}
ExecStart={{EXEC_START}}
Restart=always
RestartSec=5
User=deploy
Group=deploy
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
